# dryerMonitor.conf
[program:dryerMonitor]
directory=/opt/dryer/
command=python /opt/dryer/dryer.py
autostart=true
autorestart=true

